# pkg-scripts/postinst
#!/bin/bash

if [ $(systemctl is-enabled consul) == 'enabled' ]; then
  systemctl start consul
fi

exit 0
